# site-cookbooks/railsci/teamcity/templates/default/teamcity-agent.conf.erb
# teamcity-agent

description	"teamcity-agent"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
console output

exec su <%= node['teamcity']['user'] %> -c '/srv/TeamCity/buildAgent/bin/agent.sh run'
